# bitrix/crontab/crontab.cfg
0 6 * * * /usr/bin/php -d mbstring.func_overload=2 -d mbstring.internal_encoding=UTF-8 -f /home/sputniki/web/kawe.su/public_html/bitrix/php_interface/include/kda.exportexcel/cron_frame.php 0 >/home/sputniki/web/kawe.su/public_html/bitrix/php_interface/include/kda.exportexcel/logs/0.txt
6 0 * * * /usr/bin/php -d mbstring.func_overload=2 -d mbstring.internal_encoding=UTF-8 -f /var/www/kawe.su/data/www/kawe.su/bitrix/php_interface/include/kda.exportexcel/cron_frame.php 0 >/var/www/kawe.su/data/www/kawe.su/bitrix/php_interface/include/kda.exportexcel/logs/0.txt
